# recipes-connectivity/samba/files/nmb.service.d/override.conf
[Unit]
ConditionPathExists=/usr/sbin/nmbd

[Service]
Environment=LD_LIBRARY_PATH=/usr/lib/samba:/usr/lib/samba/private
