# recipes-connectivity/samba/files/nmb.service.d/override.conf
[Unit]
ConditionPathExists=/usr/sbin/nmbd
