# recipes-connectivity/samba/files/smb.service.d/override.conf
[Unit]
ConditionPathExists=/usr/sbin/smbd

[Service]
Environment=LD_LIBRARY_PATH=/usr/lib/samba:/usr/lib/samba/private
