# recipes-connectivity/samba/files/smb.service.d/override.conf
[Unit]
ConditionPathExists=/usr/sbin/smbd
