# recipes-support/image-portal/image-portal-service/files/image-portal-catalog-refresh.timer
[Unit]
Description=Periodic image-portal catalog refresh

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=45s
Persistent=true
Unit=image-portal-catalog-refresh.service

[Install]
WantedBy=timers.target
